# recipes-ext/automotive/neptune-ui/neptune.service
[Unit]
Description=Neptune
After=dbus.service dbus-session.service systemd-user-sessions.service

[Service]
ExecStart=/usr/bin/appcontroller /usr/bin/appman -r --dbus session -c am-config.yaml
Restart=on-failure
WorkingDirectory=/opt/neptune

[Install]
WantedBy=multi-user.target
